# output/raw_matches/Lateral Movement__wmi_remote_registry_sysmon.json
[{
  "title": "CurrentVersion Autorun Keys Modification",
  "id": "20f0ee37-5942-4e45-b7d5-c5b5db9df5cd",
  "description": "Detects modification of autostart extensibility point (ASEP) in registry.",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Microsoft-Windows-Sysmon/Operational' AND (EventID=13 AND ((TargetObject LIKE '%\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion%' ESCAPE '\\' AND (TargetObject LIKE '%\\\\ShellServiceObjectDelayLoad%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Run\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\RunOnce\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\RunOnceEx\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\RunServices\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\RunServicesOnce\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Policies\\\\System\\\\Shell%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Policies\\\\Explorer\\\\Run%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Group Policy\\\\Scripts\\\\Startup%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Group Policy\\\\Scripts\\\\Shutdown%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Group Policy\\\\Scripts\\\\Logon%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Group Policy\\\\Scripts\\\\Logoff%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Explorer\\\\ShellServiceObjects%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Explorer\\\\ShellIconOverlayIdentifiers%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Explorer\\\\ShellExecuteHooks%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Explorer\\\\SharedTaskScheduler%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Explorer\\\\Browser Helper Objects%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\PLAP Providers%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\Credential Providers%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\Credential Provider Filters%' ESCAPE '\\')) AND (NOT ((Details='(Empty)' OR TargetObject LIKE '%\\\\NgcFirst\\\\ConsecutiveSwitchCount' ESCAPE '\\' OR (Image LIKE '%\\\\AppData\\\\Local\\\\Microsoft\\\\OneDrive\\\\Update\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE '%\\\\AppData\\\\Roaming\\\\Spotify\\\\Spotify.exe' ESCAPE '\\' OR Image LIKE '%\\\\AppData\\\\Local\\\\WebEx\\\\WebexHost.exe' ESCAPE '\\') OR (Image LIKE 'C:\\\\WINDOWS\\\\system32\\\\devicecensus.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Windows\\\\system32\\\\winsat.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\Microsoft OneDrive\\\\StandaloneUpdater\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft OneDrive\\\\StandaloneUpdater\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\Microsoft OneDrive\\\\Update\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft OneDrive\\\\Update\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\Microsoft Office\\\\root\\\\integration\\\\Addons\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft Office\\\\root\\\\integration\\\\Addons\\\\OneDriveSetup.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\KeePass Password Safe 2\\\\ShInstUtil.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\Everything\\\\Everything.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft Office\\\\root\\\\integration\\\\integrator.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files\\\\Microsoft Office\\\\root\\\\integration\\\\integrator.exe' ESCAPE '\\')) OR Details IS NULL OR (Image LIKE 'C:\\\\Windows\\\\system32\\\\LogonUI.exe' ESCAPE '\\' AND (TargetObject LIKE '%\\\\Authentication\\\\Credential Providers\\\\{D6886603-9D2F-4EB2-B667-1971041FA96B}\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\Credential Providers\\\\{BEC09223-B018-416D-A0AC-523971B639F5}\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\Credential Providers\\\\{8AF662BF-65A0-4D0A-A540-A338A999D36F}\\\\%' ESCAPE '\\' OR TargetObject LIKE '%\\\\Authentication\\\\Credential Providers\\\\{27FBDB57-B613-4AF2-9D7E-4FA7A66C21AD}\\\\%' ESCAPE '\\')) OR (Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft\\\\EdgeUpdate\\\\Install\\\\%' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft\\\\EdgeWebView\\\\%' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Microsoft\\\\Edge\\\\Application\\\\msedge.exe%' ESCAPE '\\') OR Image LIKE 'C:\\\\Program Files\\\\Windows Defender\\\\MsMpEng.exe' ESCAPE '\\' OR (Image LIKE '%\\\\Microsoft\\\\Teams\\\\current\\\\Teams.exe' ESCAPE '\\' AND Details LIKE '%\\\\Microsoft\\\\Teams\\\\Update.exe --processStart %' ESCAPE '\\') OR (Image LIKE 'C:\\\\Windows\\\\system32\\\\userinit.exe' ESCAPE '\\' AND Details='ctfmon.exe /n'))) AND (NOT ((Image LIKE 'C:\\\\Windows\\\\system32\\\\regsvr32.exe' ESCAPE '\\' AND TargetObject LIKE '%DropboxExt%' ESCAPE '\\' AND Details LIKE '%A251-47B7-93E1-CDD82E34AF8B}' ESCAPE '\\') OR (TargetObject LIKE '%\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Opera Browser Assistant' ESCAPE '\\' AND Details LIKE 'C:\\\\Program Files\\\\Opera\\\\assistant\\\\browser\\_assistant.exe' ESCAPE '\\') OR (TargetObject LIKE '%\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Opera Stable' ESCAPE '\\' AND (Details LIKE 'C:\\\\Program Files\\\\Opera\\\\launcher.exe' ESCAPE '\\' OR Details LIKE 'C:\\\\Program Files (x86)\\\\Opera\\\\launcher.exe' ESCAPE '\\')) OR (TargetObject LIKE '%\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\iTunesHelper' ESCAPE '\\' AND Details LIKE '\"C:\\\\Program Files\\\\iTunes\\\\iTunesHelper.exe\"' ESCAPE '\\') OR (TargetObject LIKE '%\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\zoommsirepair' ESCAPE '\\' AND Details LIKE '\"C:\\\\Program Files\\\\Zoom\\\\bin\\\\installer.exe\" /repair' ESCAPE '\\') OR (TargetObject LIKE '%\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Greenshot' ESCAPE '\\' AND Details LIKE 'C:\\\\Program Files\\\\Greenshot\\\\Greenshot.exe' ESCAPE '\\') OR (TargetObject LIKE '%\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\GoogleDriveFS' ESCAPE '\\' AND Details LIKE 'C:\\\\Program Files\\\\Google\\\\Drive File Stream\\\\%' ESCAPE '\\' AND Details LIKE '%\\\\GoogleDriveFS.exe%' ESCAPE '\\') OR (TargetObject LIKE '%GoogleDrive%' ESCAPE '\\' AND (Details='{CFE8B367-77A7-41D7-9C90-75D16D7DC6B6}' OR Details='{A8E52322-8734-481D-A7E2-27B309EF8D56}' OR Details='{C973DA94-CBDF-4E77-81D1-E5B794FBD146}' OR Details='{51EF1569-67EE-4AD6-9646-E726C3FFC8A2}')) OR ((Details LIKE 'C:\\\\Windows\\\\system32\\\\cmd.exe /q /c rmdir /s /q \"C:\\\\Users\\\\%' ESCAPE '\\' OR Details LIKE 'C:\\\\Windows\\\\system32\\\\cmd.exe /q /c del /q \"C:\\\\Users\\\\%' ESCAPE '\\') AND Details LIKE '%\\\\AppData\\\\Local\\\\Microsoft\\\\OneDrive\\\\%' ESCAPE '\\') OR (TargetObject LIKE '%\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\{%' ESCAPE '\\' AND (Details LIKE '%\\\\AppData\\\\Local\\\\Package Cache\\\\{%' ESCAPE '\\' AND Details LIKE '%}\\\\python-%' ESCAPE '\\') AND Details LIKE '%.exe\" /burn.runonce' ESCAPE '\\') OR ((Image LIKE 'C:\\\\Program Files\\\\Common Files\\\\Microsoft Shared\\\\ClickToRun\\\\%' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\Common Files\\\\Microsoft Shared\\\\ClickToRun\\\\%' ESCAPE '\\') AND Image LIKE '%\\\\OfficeClickToRun.exe' ESCAPE '\\') OR (Image LIKE '%\\\\Microsoft\\\\Teams\\\\current\\\\Teams.exe' ESCAPE '\\' AND Details LIKE '%\\\\Microsoft\\\\Teams\\\\Update.exe --processStart%' ESCAPE '\\') OR ((Image LIKE '%C:\\\\Program Files\\\\AVG\\\\Antivirus\\\\Setup\\\\%' ESCAPE '\\' OR Image LIKE '%C:\\\\Program Files (x86)\\\\AVG\\\\Antivirus\\\\Setup\\\\%' ESCAPE '\\' OR Image LIKE '%\\\\instup.exe%' ESCAPE '\\') AND (Details LIKE '\"C:\\\\Program Files\\\\AVG\\\\Antivirus\\\\AvLaunch.exe\" /gui' ESCAPE '\\' OR Details LIKE '\"C:\\\\Program Files (x86)\\\\AVG\\\\Antivirus\\\\AvLaunch.exe\" /gui' ESCAPE '\\' OR Details='{472083B0-C522-11CF-8763-00608CC02F24}' OR Details='{472083B1-C522-11CF-8763-00608CC02F24}')) OR ((Image LIKE '%C:\\\\Program Files\\\\Avast Software\\\\Avast\\\\Setup\\\\%' ESCAPE '\\' OR Image LIKE '%C:\\\\Program Files (x86)\\\\Avast Software\\\\Avast\\\\Setup\\\\%' ESCAPE '\\' OR Image LIKE '%\\\\instup.exe%' ESCAPE '\\') AND (Details LIKE '\"C:\\\\Program Files\\\\Avast Software\\\\Avast\\\\AvLaunch.exe\" /gui' ESCAPE '\\' OR Details LIKE '\"C:\\\\Program Files (x86)\\\\Avast Software\\\\Avast\\\\AvLaunch.exe\" /gui' ESCAPE '\\')) OR ((Image LIKE 'C:\\\\Program Files\\\\AVG\\\\Antivirus\\\\avgToolsSvc.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Program Files (x86)\\\\AVG\\\\Antivirus\\\\avgToolsSvc.exe' ESCAPE '\\') AND TargetObject LIKE '%\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\StartupApproved\\\\Run\\\\%' ESCAPE '\\' AND Details='Binary Data') OR ((Image LIKE '%\\\\aurora-agent-64.exe' ESCAPE '\\' OR Image LIKE '%\\\\aurora-agent.exe' ESCAPE '\\') AND TargetObject LIKE '%\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\aurora-dashboard' ESCAPE '\\' AND Details LIKE 'C:\\\\Program Files\\\\Aurora-Agent\\\\tools\\\\aurora-dashboard.exe' ESCAPE '\\') OR (TargetObject LIKE '%\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Everything' ESCAPE '\\' AND Details LIKE '%\\\\Everything\\\\Everything.exe\" -startup' ESCAPE '\\') OR (TargetObject LIKE '%\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Discord' ESCAPE '\\' AND Details LIKE '%\\\\Discord\\\\Update.exe --processStart Discord.exe' ESCAPE '\\')))))"
  ],
  "rule_level": "medium",
  "tags": [
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1547.001"
  ],
  "count": 1,
  "matches": [
    {
      "row_id": 13,
      "Channel": "Microsoft-Windows-Sysmon/Operational",
      "Computer": "MSEDGEWIN10",
      "Details": "rundll32.exe a.dll, DllMain",
      "EventID": 13,
      "EventRecordID": 549609,
      "EventType": "SetValue",
      "Guid": "5770385F-C22A-43E0-BF4C-06F5698FFBD9",
      "Image": "C:\\Windows\\system32\\wbem\\wmiprvse.exe",
      "Keywords": "0x8000000000000000",
      "Level": 4,
      "Opcode": 0,
      "OriginalLogfile": "wmi_remote_registry_sysmon.evtx",
      "ProcessGuid": "747F3D96-041E-5FD2-0000-001024DF3B00",
      "ProcessID": 3556,
      "Provider_Name": "Microsoft-Windows-Sysmon",
      "SystemTime": "2020-12-10T11:18:54.856311Z",
      "TargetObject": "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\wmipers",
      "Task": 13,
      "ThreadID": 4972,
      "UserID": "S-1-5-18",
      "UtcTime": "2020-12-10 11:18:54.853",
      "Version": 2
    }
  ]
}]