# output/raw_matches/Lateral Movement__remote task update 4624 4702 same logonid.json
[{
  "title": "Security Eventlog Cleared",
  "id": "d99b79d2-0a6f-4f46-ad8b-260b6e17f982",
  "description": "One of the Windows Eventlogs has been cleared. e.g. caused by \"wevtutil cl\" command execution",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Security' AND ((EventID=517 AND Provider_Name='Security') OR (EventID=1102 AND Provider_Name='Microsoft-Windows-Eventlog'))"
  ],
  "rule_level": "high",
  "tags": [
    "attack.defense-evasion",
    "attack.t1070.001",
    "car.2016-04-002"
  ],
  "count": 1,
  "matches": [
    {
      "row_id": 1,
      "Channel": "Security",
      "Computer": "01566s-win16-ir.threebeesco.com",
      "EventID": 1102,
      "EventRecordID": 2171289,
      "Guid": "{fc65ddd8-d6ef-4962-83d5-6e5cfe9ce148}",
      "Keywords": "0x4020000000000000",
      "Level": 4,
      "Opcode": 0,
      "OriginalLogfile": "remote task update 4624 4702 same logonid.evtx",
      "ProcessID": 420,
      "Provider_Name": "Microsoft-Windows-Eventlog",
      "SubjectDomainName": "3B",
      "SubjectLogonId": "0x38a14",
      "SubjectUserName": "a-jbrown",
      "SubjectUserSid": "S-1-5-21-308926384-506822093-3341789130-1106",
      "SystemTime": "2020-09-02T11:47:39.499106Z",
      "Task": 104,
      "ThreadID": 996,
      "Version": 0
    }
  ]
},
{
  "title": "Suspicious Scheduled Task Update",
  "id": "614cf376-6651-47c4-9dcc-6b9527f749f4",
  "description": "Detects update to a scheduled task event that contain suspicious keywords.",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Security' AND (EventID=4702 AND (TaskContentNew LIKE '%\\\\AppData\\\\Local\\\\Temp\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\AppData\\\\Roaming\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\Users\\\\Public\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\WINDOWS\\\\Temp\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%C:\\\\Temp\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\Desktop\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\Downloads\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%\\\\Temporary Internet%' ESCAPE '\\' OR TaskContentNew LIKE '%C:\\\\ProgramData\\\\%' ESCAPE '\\' OR TaskContentNew LIKE '%C:\\\\Perflogs\\\\%' ESCAPE '\\') AND (TaskContentNew LIKE '%regsvr32%' ESCAPE '\\' OR TaskContentNew LIKE '%rundll32%' ESCAPE '\\' OR TaskContentNew LIKE '%cmd.exe</Command>%' ESCAPE '\\' OR TaskContentNew LIKE '%cmd</Command>%' ESCAPE '\\' OR TaskContentNew LIKE '%<Arguments>/c %' ESCAPE '\\' OR TaskContentNew LIKE '%<Arguments>/k %' ESCAPE '\\' OR TaskContentNew LIKE '%<Arguments>/r %' ESCAPE '\\' OR TaskContentNew LIKE '%powershell%' ESCAPE '\\' OR TaskContentNew LIKE '%pwsh%' ESCAPE '\\' OR TaskContentNew LIKE '%mshta%' ESCAPE '\\' OR TaskContentNew LIKE '%wscript%' ESCAPE '\\' OR TaskContentNew LIKE '%cscript%' ESCAPE '\\' OR TaskContentNew LIKE '%certutil%' ESCAPE '\\' OR TaskContentNew LIKE '%bitsadmin%' ESCAPE '\\' OR TaskContentNew LIKE '%bash.exe%' ESCAPE '\\' OR TaskContentNew LIKE '%bash %' ESCAPE '\\' OR TaskContentNew LIKE '%scrcons%' ESCAPE '\\' OR TaskContentNew LIKE '%wmic %' ESCAPE '\\' OR TaskContentNew LIKE '%wmic.exe%' ESCAPE '\\' OR TaskContentNew LIKE '%forfiles%' ESCAPE '\\' OR TaskContentNew LIKE '%scriptrunner%' ESCAPE '\\' OR TaskContentNew LIKE '%hh.exe%' ESCAPE '\\'))"
  ],
  "rule_level": "high",
  "tags": [
    "attack.execution",
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "count": 1,
  "matches": [
    {
      "row_id": 5,
      "ActivityID": "4F7FBBE3-7BB5-0002-EBBB-7F4FB57BD601",
      "Channel": "Security",
      "Computer": "01566s-win16-ir.threebeesco.com",
      "EventID": 4702,
      "EventRecordID": 2171293,
      "Guid": "54849625-5478-4994-A5BA-3E3B0328C30D",
      "Keywords": "0x8020000000000000",
      "Level": 0,
      "Opcode": 0,
      "OriginalLogfile": "remote task update 4624 4702 same logonid.evtx",
      "ProcessID": 632,
      "Provider_Name": "Microsoft-Windows-Security-Auditing",
      "SubjectDomainName": "3B",
      "SubjectLogonId": "0x21a8c68",
      "SubjectUserName": "a-jbrown",
      "SubjectUserSid": "S-1-5-21-308926384-506822093-3341789130-1106",
      "SystemTime": "2020-09-02T11:47:48.959767Z",
      "Task": 12804,
      "TaskContentNew": "<?xml version=\"1.0\" encoding=\"UTF-16\"?>\r\n<Task version=\"1.2\" xmlns=\"http://schemas.microsoft.com/windows/2004/02/mit/task\">\r\n  <RegistrationInfo>\r\n    <Date>2020-09-02T04:47:49.74-07:00</Date>\r\n    <Author>a-jbrown</Author>\r\n    <Description>00304d6e</Description>\r\n    <URI>\\LMST</URI>\r\n  </RegistrationInfo>\r\n  <Triggers>\r\n    <TimeTrigger>\r\n      <StartBoundary>2020-02-09T04:47:48</StartBoundary>\r\n      <EndBoundary>2020-02-09T04:47:58</EndBoundary>\r\n      <Enabled>true</Enabled>\r\n    </TimeTrigger>\r\n  </Triggers>\r\n  <Principals>\r\n    <Principal id=\"Author\">\r\n      <RunLevel>HighestAvailable</RunLevel>\r\n      <UserId>SYSTEM</UserId>\r\n    </Principal>\r\n  </Principals>\r\n  <Settings>\r\n    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>\r\n    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>\r\n    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>\r\n    <AllowHardTerminate>true</AllowHardTerminate>\r\n    <StartWhenAvailable>true</StartWhenAvailable>\r\n    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>\r\n    <IdleSettings>\r\n      <Duration>PT10M</Duration>\r\n      <WaitTimeout>PT1H</WaitTimeout>\r\n      <StopOnIdleEnd>true</StopOnIdleEnd>\r\n      <RestartOnIdle>false</RestartOnIdle>\r\n    </IdleSettings>\r\n    <AllowStartOnDemand>true</AllowStartOnDemand>\r\n    <Enabled>true</Enabled>\r\n    <Hidden>true</Hidden>\r\n    <RunOnlyIfIdle>false</RunOnlyIfIdle>\r\n    <WakeToRun>false</WakeToRun>\r\n    <ExecutionTimeLimit>PT72H</ExecutionTimeLimit>\r\n    <Priority>7</Priority>\r\n  </Settings>\r\n  <Actions Context=\"Author\">\r\n    <Exec>\r\n      <Command>cmd.exe</Command>\r\n      <Arguments>/c echo testing &gt; c:\\users\\public\\out.txt</Arguments>\r\n    </Exec>\r\n  </Actions>\r\n</Task>",
      "TaskName": "\\LMST",
      "ThreadID": 4244,
      "Version": 0
    }
  ]
},
{
  "title": "External Remote SMB Logon from Public IP",
  "id": "78d5cab4-557e-454f-9fb9-a222bd0d5edc",
  "description": "Detects successful logon from public IP address via SMB. This can indicate a publicly-exposed SMB port.",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Security' AND ((EventID=4624 AND LogonType=3) AND (NOT ((IpAddress='::1/128' OR IpAddress LIKE '10.%' ESCAPE '\\' OR IpAddress LIKE '127.%' ESCAPE '\\' OR IpAddress LIKE '172.16.%' ESCAPE '\\' OR IpAddress LIKE '172.17.%' ESCAPE '\\' OR IpAddress LIKE '172.18.%' ESCAPE '\\' OR IpAddress LIKE '172.19.%' ESCAPE '\\' OR IpAddress LIKE '172.20.%' ESCAPE '\\' OR IpAddress LIKE '172.21.%' ESCAPE '\\' OR IpAddress LIKE '172.22.%' ESCAPE '\\' OR IpAddress LIKE '172.23.%' ESCAPE '\\' OR IpAddress LIKE '172.24.%' ESCAPE '\\' OR IpAddress LIKE '172.25.%' ESCAPE '\\' OR IpAddress LIKE '172.26.%' ESCAPE '\\' OR IpAddress LIKE '172.27.%' ESCAPE '\\' OR IpAddress LIKE '172.28.%' ESCAPE '\\' OR IpAddress LIKE '172.29.%' ESCAPE '\\' OR IpAddress LIKE '172.30.%' ESCAPE '\\' OR IpAddress LIKE '172.31.%' ESCAPE '\\' OR IpAddress LIKE '192.168.%' ESCAPE '\\' OR IpAddress LIKE '169.254.%' ESCAPE '\\' OR IpAddress LIKE 'fc%' ESCAPE '\\' OR IpAddress LIKE 'fd%' ESCAPE '\\' OR IpAddress LIKE 'fe8%' ESCAPE '\\' OR IpAddress LIKE 'fe9%' ESCAPE '\\' OR IpAddress LIKE 'fea%' ESCAPE '\\' OR IpAddress LIKE 'feb%' ESCAPE '\\') OR IpAddress='-')))"
  ],
  "rule_level": "high",
  "tags": [
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.defense-evasion",
    "attack.initial-access",
    "attack.credential-access",
    "attack.t1133",
    "attack.t1078",
    "attack.t1110"
  ],
  "count": 3,
  "matches": [
    {
      "row_id": 6,
      "ActivityID": "4F7FBBE3-7BB5-0002-EBBB-7F4FB57BD601",
      "AuthenticationPackageName": "Kerberos",
      "Channel": "Security",
      "Computer": "01566s-win16-ir.threebeesco.com",
      "ElevatedToken": "%%1842",
      "EventID": 4624,
      "EventRecordID": 2171294,
      "Guid": "54849625-5478-4994-A5BA-3E3B0328C30D",
      "ImpersonationLevel": "%%1840",
      "IpAddress": "::1",
      "IpPort": "62860",
      "KeyLength": 0,
      "Keywords": "0x8020000000000000",
      "Level": 0,
      "LmPackageName": "-",
      "LogonGuid": "27FCE179-F80F-F6A6-7DF4-C247E783B072",
      "LogonProcessName": "Kerberos",
      "LogonType": 3,
      "Opcode": 0,
      "OriginalLogfile": "remote task update 4624 4702 same logonid.evtx",
      "ProcessID": 632,
      "ProcessName": "-",
      "Provider_Name": "Microsoft-Windows-Security-Auditing",
      "RestrictedAdminMode": "-",
      "SubjectDomainName": "-",
      "SubjectLogonId": "0x0",
      "SubjectUserName": "-",
      "SubjectUserSid": "S-1-0-0",
      "SystemTime": "2020-09-02T11:47:49.966623Z",
      "TargetDomainName": "THREEBEESCO.COM",
      "TargetLinkedLogonId": "0x0",
      "TargetLogonId": "0x21aa47f",
      "TargetOutboundDomainName": "-",
      "TargetOutboundUserName": "-",
      "TargetUserName": "01566S-WIN16-IR$",
      "TargetUserSid": "S-1-5-18",
      "Task": 12544,
      "ThreadID": 4244,
      "TransmittedServices": "-",
      "Version": 2,
      "VirtualAccount": "%%1843",
      "WorkstationName": "-"
    },
    {
      "row_id": 7,
      "AuthenticationPackageName": "Kerberos",
      "Channel": "Security",
      "Computer": "01566s-win16-ir.threebeesco.com",
      "ElevatedToken": "%%1842",
      "EventID": 4624,
      "EventRecordID": 2171295,
      "Guid": "54849625-5478-4994-A5BA-3E3B0328C30D",
      "ImpersonationLevel": "%%1833",
      "IpAddress": "::1",
      "IpPort": "62862",
      "KeyLength": 0,
      "Keywords": "0x8020000000000000",
      "Level": 0,
      "LmPackageName": "-",
      "LogonGuid": "860D1189-6C67-C57B-59ED-C0676A052019",
      "LogonProcessName": "Kerberos",
      "LogonType": 3,
      "Opcode": 0,
      "OriginalLogfile": "remote task update 4624 4702 same logonid.evtx",
      "ProcessID": 632,
      "ProcessName": "-",
      "Provider_Name": "Microsoft-Windows-Security-Auditing",
      "RestrictedAdminMode": "-",
      "SubjectDomainName": "-",
      "SubjectLogonId": "0x0",
      "SubjectUserName": "-",
      "SubjectUserSid": "S-1-0-0",
      "SystemTime": "2020-09-02T11:47:57.252932Z",
      "TargetDomainName": "THREEBEESCO.COM",
      "TargetLinkedLogonId": "0x0",
      "TargetLogonId": "0x21aad4a",
      "TargetOutboundDomainName": "-",
      "TargetOutboundUserName": "-",
      "TargetUserName": "01566S-WIN16-IR$",
      "TargetUserSid": "S-1-5-18",
      "Task": 12544,
      "ThreadID": 2512,
      "TransmittedServices": "-",
      "Version": 2,
      "VirtualAccount": "%%1843",
      "WorkstationName": "-"
    },
    {
      "row_id": 8,
      "AuthenticationPackageName": "Kerberos",
      "Channel": "Security",
      "Computer": "01566s-win16-ir.threebeesco.com",
      "ElevatedToken": "%%1842",
      "EventID": 4624,
      "EventRecordID": 2171296,
      "Guid": "54849625-5478-4994-A5BA-3E3B0328C30D",
      "ImpersonationLevel": "%%1833",
      "IpAddress": "::1",
      "IpPort": "62863",
      "KeyLength": 0,
      "Keywords": "0x8020000000000000",
      "Level": 0,
      "LmPackageName": "-",
      "LogonGuid": "860D1189-6C67-C57B-59ED-C0676A052019",
      "LogonProcessName": "Kerberos",
      "LogonType": 3,
      "Opcode": 0,
      "OriginalLogfile": "remote task update 4624 4702 same logonid.evtx",
      "ProcessID": 632,
      "ProcessName": "-",
      "Provider_Name": "Microsoft-Windows-Security-Auditing",
      "RestrictedAdminMode": "-",
      "SubjectDomainName": "-",
      "SubjectLogonId": "0x0",
      "SubjectUserName": "-",
      "SubjectUserSid": "S-1-0-0",
      "SystemTime": "2020-09-02T11:47:57.263194Z",
      "TargetDomainName": "THREEBEESCO.COM",
      "TargetLinkedLogonId": "0x0",
      "TargetLogonId": "0x21aadb8",
      "TargetOutboundDomainName": "-",
      "TargetOutboundUserName": "-",
      "TargetUserName": "01566S-WIN16-IR$",
      "TargetUserSid": "S-1-5-18",
      "Task": 12544,
      "ThreadID": 2504,
      "TransmittedServices": "-",
      "Version": 2,
      "VirtualAccount": "%%1843",
      "WorkstationName": "-"
    }
  ]
}]