# output/raw_matches/Privilege Escalation__sysmon_1_11_exec_as_system_via_schedtask.json
[{
  "title": "Schedule Task Creation From Env Variable Or Potentially Suspicious Path Via Schtasks.EXE",
  "id": "81325ce1-be01-4250-944f-b4789644556f",
  "description": "Detects Schtask creations that point to a suspicious folder or an environment variable often used by malware",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Microsoft-Windows-Sysmon/Operational' AND (EventID=1 AND ((((Image LIKE '%\\\\schtasks.exe' ESCAPE '\\' AND (CommandLine LIKE '% -create %' ESCAPE '\\' OR CommandLine LIKE '% /create %' ESCAPE '\\' OR CommandLine LIKE '% –create %' ESCAPE '\\' OR CommandLine LIKE '% —create %' ESCAPE '\\' OR CommandLine LIKE '% ―create %' ESCAPE '\\')) AND (CommandLine LIKE '%:\\\\Perflogs%' ESCAPE '\\' OR CommandLine LIKE '%:\\\\Users\\\\All Users\\\\%' ESCAPE '\\' OR CommandLine LIKE '%:\\\\Users\\\\Default\\\\%' ESCAPE '\\' OR CommandLine LIKE '%:\\\\Users\\\\Public%' ESCAPE '\\' OR CommandLine LIKE '%:\\\\Windows\\\\Temp%' ESCAPE '\\' OR CommandLine LIKE '%\\\\AppData\\\\Local\\\\%' ESCAPE '\\' OR CommandLine LIKE '%\\\\AppData\\\\Roaming\\\\%' ESCAPE '\\' OR CommandLine LIKE '%\\%AppData\\%%' ESCAPE '\\' OR CommandLine LIKE '%\\%Public\\%%' ESCAPE '\\')) OR (ParentCommandLine LIKE '%\\\\svchost.exe -k netsvcs -p -s Schedule' ESCAPE '\\' AND (CommandLine LIKE '%:\\\\Perflogs%' ESCAPE '\\' OR CommandLine LIKE '%:\\\\Windows\\\\Temp%' ESCAPE '\\' OR CommandLine LIKE '%\\\\Users\\\\Public%' ESCAPE '\\' OR CommandLine LIKE '%\\%Public\\%%' ESCAPE '\\'))) AND (NOT ((ParentCommandLine LIKE '%unattended.ini%' ESCAPE '\\' OR CommandLine LIKE '%update\\_task.xml%' ESCAPE '\\') OR CommandLine LIKE '%/Create /TN TVInstallRestore /TR%' ESCAPE '\\' OR (CommandLine LIKE '%/Create /Xml %' ESCAPE '\\' AND CommandLine LIKE '%\\\\Temp\\\\.CR.%' ESCAPE '\\' AND CommandLine LIKE '%\\\\Avira\\_Security\\_Installation.xml%' ESCAPE '\\') OR ((CommandLine LIKE '%/Create /F /TN%' ESCAPE '\\' AND CommandLine LIKE '%/Xml %' ESCAPE '\\' AND CommandLine LIKE '%\\\\Temp\\\\%' ESCAPE '\\' AND CommandLine LIKE '%Avira\\_%' ESCAPE '\\') AND (CommandLine LIKE '%.tmp\\\\UpdateFallbackTask.xml%' ESCAPE '\\' OR CommandLine LIKE '%.tmp\\\\WatchdogServiceControlManagerTimeout.xml%' ESCAPE '\\' OR CommandLine LIKE '%.tmp\\\\SystrayAutostart.xml%' ESCAPE '\\' OR CommandLine LIKE '%.tmp\\\\MaintenanceTask.xml%' ESCAPE '\\')) OR (CommandLine LIKE '%\\\\Temp\\\\%' ESCAPE '\\' AND CommandLine LIKE '%/Create /TN \"klcp\\_update\" /XML %' ESCAPE '\\' AND CommandLine LIKE '%\\\\klcp\\_update\\_task.xml%' ESCAPE '\\')))))"
  ],
  "rule_level": "medium",
  "tags": [
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.execution",
    "attack.t1053.005"
  ],
  "count": 1,
  "matches": [
    {
      "row_id": 2,
      "Channel": "Microsoft-Windows-Sysmon/Operational",
      "CommandLine": "\"C:\\Windows\\System32\\schtasks.exe\" /create /xml c:\\users\\ieuser\\appdata\\local\\temp\\elevator.xml /tn elevator",
      "Company": "Microsoft Corporation",
      "Computer": "IEWIN7",
      "CurrentDirectory": "c:\\Users\\IEUser\\Downloads\\WinPwnage-master\\WinPwnage-master\\",
      "Description": "Manages scheduled tasks",
      "EventID": 1,
      "EventRecordID": 16243,
      "FileVersion": "6.1.7600.16385 (win7_rtm.090713-1255)",
      "Guid": "5770385F-C22A-43E0-BF4C-06F5698FFBD9",
      "Hashes": "SHA1=8A7E8B05A122B768AB85466B2A3DAF7A358F90F4,MD5=2003E9B15E1C502B146DAD2E383AC1E3,SHA256=15018D0093BEFABBA8B927743191030D1F8C17BB97FDB48C2FC3EAB20E2D4B3D,IMPHASH=D92C80D49382091310FB8DB089F856A9",
      "IMPHASH": "D92C80D49382091310FB8DB089F856A9",
      "Image": "C:\\Windows\\System32\\schtasks.exe",
      "IntegrityLevel": "High",
      "Keywords": "0x8000000000000000",
      "Level": 4,
      "LogonGuid": "365ABB72-5DEC-5CD7-0000-00204A380100",
      "LogonId": "0x1384a",
      "MD5": "2003E9B15E1C502B146DAD2E383AC1E3",
      "Opcode": 0,
      "OriginalLogfile": "sysmon_1_11_exec_as_system_via_schedtask.evtx",
      "ParentCommandLine": "python  winpwnage.py -u elevate -i 4 -p c:\\Windows\\System32\\cmd.exe",
      "ParentImage": "C:\\Python27\\python.exe",
      "ParentProcessGuid": "365ABB72-6998-5CD7-0000-00104E422200",
      "ParentProcessId": 2740,
      "ProcessGuid": "365ABB72-699E-5CD7-0000-001073582200",
      "ProcessID": 1996,
      "Product": "Microsoft® Windows® Operating System",
      "Provider_Name": "Microsoft-Windows-Sysmon",
      "SHA1": "8A7E8B05A122B768AB85466B2A3DAF7A358F90F4",
      "SHA256": "15018D0093BEFABBA8B927743191030D1F8C17BB97FDB48C2FC3EAB20E2D4B3D",
      "SystemTime": "2019-05-12T00:32:30.211387Z",
      "Task": 1,
      "TerminalSessionId": 1,
      "ThreadID": 1832,
      "User": "IEWIN7\\IEUser",
      "UserID": "S-1-5-18",
      "UtcTime": "2019-05-12 00:32:30.023",
      "Version": 5
    }
  ]
},
{
  "title": "Scheduled Task Creation Via Schtasks.EXE",
  "id": "92626ddd-662c-49e3-ac59-f6535f12d189",
  "description": "Detects the creation of scheduled tasks by user accounts via the \"schtasks\" utility.",
  "sigmafile": "",
  "sigma": [
    "SELECT * FROM logs WHERE Channel='Microsoft-Windows-Sysmon/Operational' AND (EventID=1 AND ((Image LIKE '%\\\\schtasks.exe' ESCAPE '\\' AND CommandLine LIKE '% /create %' ESCAPE '\\') AND (NOT (User LIKE '%AUTHORI%' ESCAPE '\\' OR User LIKE '%AUTORI%' ESCAPE '\\')) AND (NOT ((ParentImage LIKE 'C:\\\\Program Files\\\\Microsoft Office\\\\root\\\\integration\\\\integrator.exe' ESCAPE '\\' OR ParentImage LIKE 'C:\\\\Program Files (x86)\\\\Microsoft Office\\\\root\\\\integration\\\\integrator.exe' ESCAPE '\\') AND (Image LIKE 'C:\\\\Windows\\\\System32\\\\schtasks.exe' ESCAPE '\\' OR Image LIKE 'C:\\\\Windows\\\\SysWOW64\\\\schtasks.exe' ESCAPE '\\') AND CommandLine LIKE '%Microsoft\\\\Office\\\\Office Performance Monitor%' ESCAPE '\\'))))"
  ],
  "rule_level": "low",
  "tags": [
    "attack.execution",
    "attack.persistence",
    "attack.privilege-escalation",
    "attack.t1053.005",
    "attack.s0111",
    "car.2013-08-001",
    "stp.1u"
  ],
  "count": 1,
  "matches": [
    {
      "row_id": 2,
      "Channel": "Microsoft-Windows-Sysmon/Operational",
      "CommandLine": "\"C:\\Windows\\System32\\schtasks.exe\" /create /xml c:\\users\\ieuser\\appdata\\local\\temp\\elevator.xml /tn elevator",
      "Company": "Microsoft Corporation",
      "Computer": "IEWIN7",
      "CurrentDirectory": "c:\\Users\\IEUser\\Downloads\\WinPwnage-master\\WinPwnage-master\\",
      "Description": "Manages scheduled tasks",
      "EventID": 1,
      "EventRecordID": 16243,
      "FileVersion": "6.1.7600.16385 (win7_rtm.090713-1255)",
      "Guid": "5770385F-C22A-43E0-BF4C-06F5698FFBD9",
      "Hashes": "SHA1=8A7E8B05A122B768AB85466B2A3DAF7A358F90F4,MD5=2003E9B15E1C502B146DAD2E383AC1E3,SHA256=15018D0093BEFABBA8B927743191030D1F8C17BB97FDB48C2FC3EAB20E2D4B3D,IMPHASH=D92C80D49382091310FB8DB089F856A9",
      "IMPHASH": "D92C80D49382091310FB8DB089F856A9",
      "Image": "C:\\Windows\\System32\\schtasks.exe",
      "IntegrityLevel": "High",
      "Keywords": "0x8000000000000000",
      "Level": 4,
      "LogonGuid": "365ABB72-5DEC-5CD7-0000-00204A380100",
      "LogonId": "0x1384a",
      "MD5": "2003E9B15E1C502B146DAD2E383AC1E3",
      "Opcode": 0,
      "OriginalLogfile": "sysmon_1_11_exec_as_system_via_schedtask.evtx",
      "ParentCommandLine": "python  winpwnage.py -u elevate -i 4 -p c:\\Windows\\System32\\cmd.exe",
      "ParentImage": "C:\\Python27\\python.exe",
      "ParentProcessGuid": "365ABB72-6998-5CD7-0000-00104E422200",
      "ParentProcessId": 2740,
      "ProcessGuid": "365ABB72-699E-5CD7-0000-001073582200",
      "ProcessID": 1996,
      "Product": "Microsoft® Windows® Operating System",
      "Provider_Name": "Microsoft-Windows-Sysmon",
      "SHA1": "8A7E8B05A122B768AB85466B2A3DAF7A358F90F4",
      "SHA256": "15018D0093BEFABBA8B927743191030D1F8C17BB97FDB48C2FC3EAB20E2D4B3D",
      "SystemTime": "2019-05-12T00:32:30.211387Z",
      "Task": 1,
      "TerminalSessionId": 1,
      "ThreadID": 1832,
      "User": "IEWIN7\\IEUser",
      "UserID": "S-1-5-18",
      "UtcTime": "2019-05-12 00:32:30.023",
      "Version": 5
    }
  ]
}]